# etc/udev/rules.d/80-wifi-start.rules
SUBSYSTEM=="net", ACTION=="add", DRIVERS=="?*", KERNEL=="wlan*", RUN+="/etc/network/wlan-runner $env{INTERFACE} up"
SUBSYSTEM=="net", ACTION=="remove", DRIVERS=="?*", KERNEL=="wlan*", RUN+="/etc/network/wlan-runner $env{INTERFACE} down"
